# deploy/cafe24-easyup/crm-site2.service
[Unit]
Description=CRM Site2 Node.js Service
After=network.target

[Service]
Type=simple
User=crm
Group=crm
WorkingDirectory=/opt/crm/site2
EnvironmentFile=/opt/crm/site2/.env.production
ExecStart=/usr/bin/node /opt/crm/site2/dist/index.js
Restart=always
RestartSec=3
StandardOutput=append:/var/log/crm-site2.out.log
StandardError=append:/var/log/crm-site2.err.log
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
